# scripts/com.rina.voicevox-engine.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.rina.voicevox-engine</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/ocmm/rina/vendor/voicevox/run</string>
        <string>--host</string>
        <string>127.0.0.1</string>
        <string>--port</string>
        <string>50021</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/ocmm/rina/logs/voicevox-engine.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/ocmm/rina/logs/voicevox-engine.log</string>
</dict>
</plist>
